# angrybeta.service
[Unit]
Description=git repository service
After=network.target

[Service]
ExecStart="/usr/local/bin/angrybeta"

[Install]
WantedBy=default.target